# quadlets/n8n.container
[Unit]
Description=traefik container
Requires=soc-network.service
Requires=postgres.service
Requires=traefik.service
After=postgres.service

[Container]
ContainerName=n8n
Image=docker.n8n.io/n8nio/n8n:latest
Network=soc.network
Volume=n8n.volume:/home/node/.n8n:Z
Environment=GENERIC_TIMEZONE=Asia/Singapore
Environment=TZ=Asia/Singapore
Environment=N8N_ENFORCE_SETTINGS_FILE_PERMISSIONS=true
Environment=N8N_RUNNERS_ENABLED=true
Environment=DB_TYPE=postgresdb
Environment=DB_POSTGRESDB_DATABASE=n8n
Environment=DB_POSTGRESDB_HOST=postgres
Environment=DB_POSTGRESDB_PORT=5432
Environment=DB_POSTGRESDB_USER=n8n
Environment=DB_POSTGRESDB_PASSWORD=password
Environment=DB_POSTGRESDB_SCHEMA=public
Label=traefik.enable=true
Label=traefik.http.routers.n8n.rule=Host(`n8n.vx`)
Label=traefik.http.routers.n8n.entrypoints=websecure
Label=traefik.http.routers.dashboard.tls=true

[Install]
WantedBy=multi-user.target